# dotfiles/systemd/system/getty@tty2.service.d/override.conf
[Service]
ExecStart=
ExecStart=agetty --nohostname -JNno aleks %I
