# systemd/system/getty@.service.d/override.conf
[Service]
ExecStart=
ExecStart=agetty -Nna <user> %I
