# systemd/system/getty@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/agetty -Nna <user> %I
