# 10-udev.rules
SUBSYSTEM=="drm", ACTION=="change", ENV{DISPLAY}=":1", ENV{XAUTHORITY}="<home>/.local/share/sx/xauthority", RUN+="<home>/.local/bin/monitor"
