# autostart/xbanish.desktop
[Desktop Entry]
Name=xbanish
Exec=xbanish
Type=Application
