# autostart/xkbcomp.desktop
[Desktop Entry]
Name=xkbcomp
Exec=sh -c 'xkbcomp "$XDG_CONFIG_HOME"/xkeymap "$DISPLAY"'
Type=Application
