# autostart/xkbcomp.desktop
[Desktop Entry]
Name=xkbcomp
Exec=sh -c "xkbcomp "$HOME"/.config/xkb "$DISPLAY""
Type=Application
