# autostart/xkbmap.desktop
[Desktop Entry]
Name=xkbmap
Exec=xkbmap
Type=Application
