# dotfiles/systemd/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=agetty -JNina aleks %I
