# systemd/system/getty@tty2.service.d/override.conf
[Service]
ExecStart=
ExecStart=agetty --nohostname -Nno aleks %I
