# systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=agetty -JNna aleks %I
